# content/etc/s6-overlay/s6-rc.d/dcc/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo '** Starting dccifd **'
set -f
exec s6-setuidgid dcc /usr/libexec/dccifd ${DCCIFD_ARGS}